# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid USERNAME
chdir /home/USERNAME/sites/SITE_DOMAIN/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    SITE.wsgi:application
